# config/supervisord_nbviewer.conf
[program:nbviewer]
command=nbconvert_flask.sh
autostart=false
