# config/supervisord_redis.conf
[program:redis]
command=konix_redis-server.sh
