# config/supervisord_mpd.conf
[program:mpd]
command=mpd --no-daemon "%(ENV_KONIX_MPD_CONF)s"
