# config/supervisord_nbviewer.conf
[program:nbviewer]
command=nbconvert_flask.sh --hide-input
autostart=false
