# config/supervisord_mopidy.conf
[program:mopidy]
command=mopidy -v
priority=1000