# config/supervisord_notebook.conf
[program:notebook]
command=konix_notebook.sh
autostart=false
