# config/supervisord_gpg_agent.conf
[program:gpg-agent]
command=konix_gpg_agent.sh
